# downloader.service
[Unit]
Description=YouTube Downloader Web App
After=network.target

[Service]
Type=simple
User=raspberry
WorkingDirectory=/home/raspberry/Programs/download_all
ExecStart=/home/raspberry/Programs/download_all/venv/bin/python downloader_app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
